# elemento-rename-ifaces.service
[Unit]
Description=elemento-rename-ifaces

[Service]
Type=oneshot
ExecStart=/usr/local/bin/elemento-rename-if.py

[Install]
WantedBy=network.target

